# files/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/etc/init.d/nginx start

[program:sks]
command=/etc/init.d/sks start
